# ssh-copilot-bot.service
[Unit]
Description=SSH Copilot Bot Service
After=network.target

[Service]
User=your_user         # replace with the system user that will run the bot
Group=your_group       # replace with the group for the bot user
WorkingDirectory=/path/to/ssh-copilot-bot
ExecStart=/path/to/ssh-copilot-bot/venv/bin/python /path/to/ssh-copilot-bot/infra-bot.py
Restart=on-failure
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target

